# mt5/vps_infra/systemd/mt5b3-ticks.timer
[Unit]
Description=Timer 18:50 BRT dias uteis - coleta ticks B3

[Timer]
OnCalendar=Mon..Fri 18:50 America/Sao_Paulo
Persistent=true

[Install]
WantedBy=timers.target
